# bash_script/old_file_destroyer/old_file_destroyer.timer
[Unit]
Description=Run file_destroyer every two minutes
[Timer]
OnCalendar=*:0/2
Persistent=true
[Install]
WantedBy=timers.target
